# alpine-consul-nomad/root/etc/services.d/nomad/run
#!/usr/bin/with-contenv sh

# Bind to the external (LAN) IP
EXTERNAL_IP=`awk 'NR==1 {print $1}' /etc/hosts`

exec s6-setuidgid nomad nomad agent -bind=${EXTERNAL_IP} -config /etc/nomad/conf.d/;
